# infra/deploy/systemd/voyagent-pg-backup-offsite.timer
[Unit]
Description=Nightly trigger for voyagent-pg-backup-offsite.service
Documentation=https://github.com/voyagent/gbs-agentic-travel/blob/main/docs/BACKUPS.md

[Timer]
# 03:30 UTC every day — 90 min after the local backup at 02:00 UTC so the
# freshest local dump is always present before we upload.
OnCalendar=*-*-* 03:30:00 UTC
Persistent=true
AccuracySec=1min
Unit=voyagent-pg-backup-offsite.service

[Install]
WantedBy=timers.target
