# infra/deploy/staging/systemd-units/voyagent-web-staging.service
[Unit]
Description=Voyagent Web (staging)
Documentation=https://github.com/voyagent/gbs-agentic-travel/blob/main/docs/STAGING.md
After=network-online.target voyagent-api-staging.service
Wants=network-online.target

[Service]
Type=simple
User=voyagent
Group=voyagent
WorkingDirectory=/opt/voyagent-staging/repo/apps/web
EnvironmentFile=/opt/voyagent-staging/.env.staging
# Staging web listens on 3021 (prod uses 3011).
Environment=PORT=3021
Environment=NODE_ENV=production
Environment=VOYAGENT_ENV=staging
ExecStart=/usr/bin/env pnpm start
Restart=on-failure
RestartSec=5s
KillSignal=SIGTERM
TimeoutStopSec=30s
StandardOutput=journal
StandardError=journal
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/voyagent-staging/logs /opt/voyagent-staging/repo

[Install]
WantedBy=multi-user.target
